# bot.service
[Unit]
Description=Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/root/telegram-bot
ExecStart=/usr/bin/python3 /root/telegram-bot/bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Ограничения ресурсов
MemoryLimit=256M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
